# backend/systemd/psychic-homily-stage.service
[Unit]
Description=Psychic Homily Backend Stage
After=network.target docker.service
Requires=docker.service
Wants=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/psychic-homily-stage
ExecStart=/opt/psychic-homily-stage/psychic-homily-stage
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=psychic-homily-stage

# Environment variables
Environment=ENVIRONMENT=stage
Environment=API_ADDR=0.0.0.0:8081
EnvironmentFile=/opt/psychic-homily-stage/backend/.env.stage

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/psychic-homily-stage

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
